# roles/cabify/templates/upstart/cabify.conf.j2
description "Cabify - Example App"

start on (started networking)
stop on runlevel [!2345]

respawn

script
    python {{ cabify_path }}/cabify.py
end script
